# docker/crontab
0 0 * * * cd /prod/; /opt/pysetup/.venv/bin/python -m catpick >> /var/log/cron.log 2>&1
